# ansible_role/roles/blog_role/templates/mdbook/mdbook_systemd.service.j2
[Unit]
Description=systemd book
AssertPathIsDirectory=/blog/books/systemd

[Service]
#MemoryLimit=250M
MemoryHigh=250M
CPUQuota=70%
WorkingDirectory=/blog/books/systemd
ExecStartPre=bash -c '/root/.cargo/bin/mdbook clean'
ExecStart=/root/.cargo/bin/mdbook serve --open -n 0.0.0.0 -p 8881
ExecStartPost=bash -c 'sleep 1;/usr/bin/cp mytheme/fonts/Cascadia.woff book/mytheme/fonts'
#Restart=always
#RestartSec=15s

[Install]
WantedBy=multi-user.target
